# util/repairIface.service
[Unit]
Description=Repair the routes for the second network interface
After=network.target

[Service]
ExecStart=/root/repairIface.sh 

[Install]
WantedBy=multi-user.target
