# roles/lb-routes/files/add-routes-lb.service
[Unit]
After=repairRoutes.service

[Service]
ExecStart=/etc/kubernetes/add-routes-lb.sh

[Install]
WantedBy=multi-user.target
