# util/sed.src/repairRoutes.service
[Unit]
Description=Repair routes on multi NIC setup
After=network.target

[Service]
ExecStart=/root/repairRoutes.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
