# roles/network_ovn/templates/ovn-controller.service.j2
[Unit]
Description=Open virtual net controller
After=openvswitch-switch.service

[Service]
ExecStart=/usr/share/openvswitch/scripts/ovn-ctl start_controller

[Install]
WantedBy=default.target
